# projects/json/resources/json.service
[Unit]
Description=Json Service
After=network.target

[Service]
ExecStart=python3 /opt/kfishmonger/projects/json/start.py

[Install]
WantedBy=graphical.target